# xrootd/config/default/etc/supervisord.d/10-xrootd.conf
[program:startup]
command=/opt/cache-update.sh
priority=-1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
autorestart=false

[program:fetch-crl]
command=/opt/fetch-crl.sh
autorestart=true

[program:clustered-cmsd]
command=/opt/start-cmsd.sh
user=xrootd
autorestart=true

[program:clustered-xrootd]
command=/opt/start-xrootd.sh
autorestart=true
# These provide xrootd with the ability to switch UIDs/GIDs for reading/writing files.
CapabilityBoundingSet=CAP_SETUID CAP_SETGID CAP_DAC_OVERRIDE
Capabilities=CAP_SETGID+p CAP_SETUID+p

[program:xrootd-mon]
command=python3 /opt/xroot-mon.py
autorestart=true

[program:gunicorn]
command=bash -c "cd /opt; gunicorn --workers=2  http-server:app -b 0.0.0.0:8000"
autorestart=true

[program:node-exporter]
command=node_exporter --collector.netdev.address-info
autorestart=true
