# xrootd/config/default/etc/supervisord.d/11-usergroupchecker.conf
[program:usergroupchecker]
command=python3 /opt/scripts/usergroupchecker.py
autorestart=true
